# deployResources/gateway-tower-heart.ini
[program:gateway-tower-heart]
directory = /home/yyds/gateway/heart
command = /usr/local/jdk1.8.0_121/bin/java -Dspring.profiles.active=dev -jar -Xms128m -Xmx512m -server gateway-tower-heart-1.0.0.jar
priority = 125
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20
stdout_logfile = /home/yyds/gateway/heart/gateway-tower-heart.log