# deployResources/scloud-server-history.ini
[program:scloud-server-history]
directory = /home/yyds/scloud/history
command = /usr/local/jdk1.8.0_121/bin/java -Dspring.profiles.active=dev -jar -Xms128m -Xmx512m -server scloud-server-history-1.0.0.jar
priority = 120
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 20
stdout_logfile = /home/yyds/scloud/history/scloud-server-history.log